# root/etc/s6-overlay/s6-rc.d/svc-weewx/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

cd /home/abc

exec /usr/local/bin/weewxd
