# root/etc/s6-overlay/s6-rc.d/svc-weewx/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

cd /etc/weewx || exit 1

exec /lsiopy/bin/weewxd
