# root/etc/s6-overlay/s6-rc.d/init-weewx-config/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

if [[ -f "/etc/weewx/weewx.conf" ]]; then
    exit 0
fi

cp -a /etc/dist-weewx/. /etc/weewx/
